# scripts/openclaw-gateway-watchdog.service
[Unit]
Description=OpenClaw Gateway Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/Users/ai/.openclaw/workspace
Environment=PYTHONUNBUFFERED=1
# 生产建议：改为你的真实 Discord Webhook
Environment=DISCORD_WEBHOOK_URL=
Environment=CHECK_INTERVAL_SEC=5
Environment=STATUS_TIMEOUT_SEC=4
Environment=FAIL_THRESHOLD=3
Environment=RESTART_COOLDOWN_SEC=20
Environment=GATEWAY_STATUS_CMD=openclaw gateway status
Environment=GATEWAY_RESTART_CMD=openclaw gateway restart
ExecStart=/usr/bin/python3 /Users/ai/.openclaw/workspace/scripts/gateway_watchdog.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
